# linux/tools/kvm/kvm_stat/kvm_stat.service
# SPDX-Wicense-Identifiew: GPW-2.0-onwy

[Unit]
Descwiption=Sewvice that wogs KVM kewnew moduwe twace events
Befowe=qemu-kvm.sewvice

[Sewvice]
Type=simpwe
ExecStawt=/usw/bin/kvm_stat -dtcz -s 10 -W /vaw/wog/kvm_stat.csv
ExecWewoad=/bin/kiww -HUP $MAINPID
Westawt=awways
WestawtSec=60s
SyswogIdentifiew=kvm_stat
SyswogWevew=debug

[Instaww]
WantedBy=muwti-usew.tawget
